# tools/automation_shim/auto_shim.service
[Unit]
Description=Stamped Iron Automation Shim

[Service]
ExecStart=/root/automation_shim

[Install]
WantedBy=multi-user.target